# riak/supervisor.conf
[supervisord]
nodaemon=true

[program:riak]
command=riak start
pidfile=/var/log/riak/riak.pid
autostart=true
autorestart=true

